# cron/locate.cron
0 */4 * * *     /usr/bin/updatedb -U /home/ben
